# deploy/systemd/codify-ip-geoip-update.service
[Unit]
Description=Refresh codifyworx IP DB-IP Lite geolocation databases
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/app/ip
ExecStart=/app/ip/scripts/refresh-dbip-lite.sh
